# PiSource/install/terrarium.service
[Unit]
Description=Terrarium Controller Service
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=on-failure
RestartSec=10
User=terrarium
WorkingDirectory=/opt/terrarium
Environment="ASPNETCORE_URLS=http://0.0.0.0:5000"
Environment="ASPNETCORE_ENVIRONMENT=Production"
ExecStart=/opt/terrarium/run.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
