# PiSource/install/terrarium.service
[Unit]
Description=Terrarium Controller Service
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=on-failure
RestartSec=10
User=terrarium
WorkingDirectory=/opt/terrarium
Environment="ASPNETCORE_URLS=http://0.0.0.0:5000"
Environment="ASPNETCORE_ENVIRONMENT=Production"
ExecStart=/opt/terrarium/run.sh
ExecStartPost=/bin/sh -c 'sleep 2; mkdir -p /opt/terrarium/logs; mjpg_streamer -i "input_raspicam.so -fps 15 -w 1280 -h 720" -o "output_http.so -p 8080 -w /usr/local/share/mjpg-streamer/www" > /opt/terrarium/logs/mjpg-streamer.log 2>&1 &'
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
